# service/com.lobbyha.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.lobbyha.server</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>--import</string>
    <string>tsx</string>
    <string>packages/server/src/server.ts</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/usr/local/opt/lobbyha</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key>
    <string>production</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/usr/local/var/log/lobbyha.log</string>

  <key>StandardErrorPath</key>
  <string>/usr/local/var/log/lobbyha.log</string>
</dict>
</plist>
